# docker/supervisord/server/services/0.init-script.conf
[program:init_script]
command=/bin/bash -c "chmod 764 /scripts/init-server.sh && /scripts/init-server.sh"
user=root
chmod = 0744
numprocs=1
autorestart=false
startretries=1
