# docker/supervisord/server/services/1.php-fpm.conf
[program:php-fpm]
command=/bin/bash -c '/scripts/run-command.sh "php-fpm7.1 -R -F"'
autostart=true
autorestart=true
stopasgroup=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true