# docker/supervisord/server/services/2.rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
user=root
numprocs=1
autostart=true
autorestart=true
