# docker/supervisord/server/services/2.rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd
user=root
autorestart=false
