# docker/supervisord/nginx/services/1.nginx.conf
[program:nginx]
command=pidproxy /var/run/nginx.pid /bin/sh -c "chmod 764 /scripts/init-nginx.sh && /scripts/init-nginx.sh"
user=root
chmod = 0744
autorestart=false
autostart=true
priority=5
