# docker/supervisord/nginx/services/1.nginx.conf
[program:nginx]
command=pidproxy /var/run/nginx.pid /bin/sh -c "/scripts/init-nginx.sh"
autorestart=true
autostart=true
priority=5
